# zwave-js-ui/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bashio
bashio::log.warning "nginx exited, code: ${1}"
