# zwave-js-ui/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
exec nginx
