# rootfs/etc/s6-overlay/s6-rc.d/zwave-js-ui/finish
#!/command/with-contenv bashio
bashio::log.warning "Z-Wave JS UI exited, code: ${1}"
