# rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
tempio -conf /data/options.json -template /etc/nginx/templates/ingress.gtpl -out /etc/nginx/servers/ingress.conf
