# rootfs/etc/s6-overlay/s6-rc.d/discovery/run
#!/command/with-contenv bashio
bashio::net.wait_for 3000 localhost 1800
sleep 10
payload=$(bashio::var.json host "$(hostname)" port "^3000")
if bashio::discovery "zwave_js" "${payload}" > /dev/null; then
  bashio::log.info "Discovery info sent to Home Assistant."
else
  bashio::log.error "Discovery message failed!"
fi
